# docs/systemd/opensvc/osvc-ingress-sync.service
[Unit]
Description=Generate Traefik failover config from OpenSVC + Docker
Wants=network-online.target
After=network-online.target docker.service

[Service]
Type=oneshot
WorkingDirectory=/home/ubuntu/my-media-stack
ExecStart=/home/ubuntu/my-media-stack/scripts/osvc_ingress_sync.sh


