# docs/systemd/opensvc/osvc-l4-sync.timer
[Unit]
Description=Periodically regenerate HAProxy L4 config (OpenSVC + Docker)

[Timer]
OnBootSec=30s
OnUnitActiveSec=30s
AccuracySec=5s

[Install]
WantedBy=timers.target


