# vendor/docker/openvpn/vpn-failover.service
[Unit]
Description=VPN Failover Service
After=network.target docker.service
Wants=network.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/bin/python3 /etc/vpn-failover/vpn-failover-service.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log /var/run /etc/openvpn

# Environment
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
